# root/app/s6/musicbot/run
#!/bin/sh
cd "$ENV_MUSICBOT_DIR"
venv python3 "$ENV_MUSICBOT_DIR/run.py"
